# content/etc/profile.d/steam-display-export.sh
# naïvely set DISPLAY so that we can load Steam in our hosted X session.
export DISPLAY=:0
